# fileserver/projects/multi/uwsgi_upstart.conf
{% set python_path = project_data.python_path|default('') %}
{% set uwsgi_ini = project_data.uwsgi_ini|default('server/uwsgi.ini') %}
{% set django_project_name = project_data.django_project_name|default(project_name) %}
{% set wsgi_module = project_data.wsgi_module|default(django_project_name+'.wsgi') %}
{% set uwsgi_bin = project_data.uwsgi_bin|default('/usr/local/bin/uwsgi') %}
#######################################
# Documentation
#######################################

description "{{ project_name }}"



#######################################
# start/stop events
#######################################

start on runlevel [2345]
stop on runlevel [06]


#######################################
# Prerequisite /var/run directories
# one for socket
# other for logrotate signal file
#######################################

pre-start script
  # ensure run directory for sockets
  run_dir=/var/run/{{ project_name }}

  mkdir -p $run_dir
  chown www-data:www-data $run_dir
  chmod 0775 $run_dir

end script


#######################################
# main startup
#
# environment specific settings are here
# project specific settings are in the ini file
#######################################
script

  command="{{ uwsgi_bin }}

  --uid www-data
  --gid www-data
  --socket /var/run/{{ project_name }}/{{ project_name }}.sock
  --vacuum

  --virtualenv /opt/virtualenvs/{{ project_name }}
  --pythonpath /opt/src/{{ project_name }}/{{ python_path }}
  --module={{ wsgi_module }}:application

  --logfile-chmod 664
  --logto2 /var/log/{{ project_name }}/uwsgi.log

  --ini /opt/src/{{ project_name }}/{{ uwsgi_ini }}
  --enable-threads
  --single-interpreter
  --die-on-term

  {% if project_data.db_name|default() %}
  --env=DB_NAME={{ project_data.db_name }}
  --env=DB_PASS={{ project_data.db_pass }}
  --env=DB_USER={{ project_data.db_user }}
  --env=DB_HOST={{ project_data.db_host }}
  {% endif %}

  "

  exec $command


end script
